# sources/queued.service
[Unit]
Description=Queue daemon

[Service]
Type=forking
ExecStart=/usr/bin/queued
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
